# systemd/growattserver.socket
[Unit]
Description=Growatt Server Socket

[Socket]
ListenStream=5279
Accept=yes

[Install]
WantedBy=sockets.target
