# systemd/growattproxy.socket
[Unit]
Description=Growatt Proxy Socket

[Socket]
ListenStream=5279
Accept=yes

[Install]
WantedBy=sockets.target
